# deploy/playbooks/templates/supervisor-celery.conf.j2
[program:{{ app_name }}-celery]
command={{ app_dir }}/venv/bin/celery -A celery_app.celery_app worker --loglevel=info --concurrency=4
directory={{ app_dir }}/src
user={{ app_user }}
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/{{ app_name }}-celery.log
stderr_logfile=/var/log/supervisor/{{ app_name }}-celery_error.log
environment=PATH="{{ app_dir }}/venv/bin:%(ENV_PATH)s",
            PYTHONPATH="{{ app_dir }}/src:{{ app_dir }}",
            DATABASE_URL="postgresql://{{ db_user | default('linkedclaims') }}:{{ db_password }}@localhost/{{ db_name | default('linkedclaims_extraction') }}",
            REDIS_HOST="{{ redis_host | default('localhost') }}",
            REDIS_PORT="{{ redis_port | default('6379') }}",
            C_FORCE_ROOT="true"